# .config/systemd/user/rand-bg.timer
[Unit]
Description=Run Random Background Every 45 Minutes

[Timer]
OnStartupSec=0
OnUnitActiveSec=45min

[Install]
WantedBy=timers.target
