# .config/systemd/user/rand-bg.service
[Unit]
Description=Random Background Changer

[Service]
Type=oneshot
ExecStart=/home/nasreddin/.local/bin/rand-bg

[Install]
WantedBy=default.target
