# .config/systemd/user/vdirsyncer.service
[Unit]
Description=Sync vdirsyncer calendars

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
