# .config/autostart/OpenRGB.desktop
[Desktop Entry]
Type=Application
Encoding=UTF-8
Name=OpenRGB
Comment=Control RGB lighting
Exec=openrgb -p nas
Icon=OpenRGB
Terminal=false
Categories=Utility;
