# .config/autostart/unclutter.desktop
[Desktop Entry]
Name=Unclutter
Exec=unclutter
Terminal=false
Type=Application
